# port_proxy/proxy_port/config/ssh_config_template.conf
[program:$__task_name__]
directory=$__directory__
command=$__command__
priority=999
startsecs=5
autorestart=true
startretries=3
stdout_logfile=$__logfile__
stopasgroup=true

